# roles/keepassxc/files/systemd/user/keepassxc.service
[Unit]
Description=Start keepass

[Service]
Type=simple
ExecStart=/usr/bin/keepassxc

[Install]
WantedBy=graphical-session.target
