# roles/mako/files/systemd/user/mako.service
[Unit]
Description=Mako notification daemon

[Service]
Type=simple
ExecStart=/usr/bin/mako

[Install]
WantedBy=graphical-session.target
