# roles/filesystem/files/systemd/user/udiskie.service
[Unit]
Description=Manage usb automount

[Service]
Type=simple
ExecStart=/usr/bin/udiskie -T

[Install]
WantedBy=graphical-session.target
