# roles/network/files/systemd/user/nm-applet.service
[Unit]
Description=network manager systemtray

[Service]
Type=simple
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=graphical-session.target
